# supervisord.conf
[supervisord]
nodaemon=true

[program:mysql]
command=service mysql start
