# supervisord.conf
[supervisord]
nodaemon=true

[program:mysql]
command=service mysql start


[program:runserver]
directory=/opt/dmoj/site
priority=900
command='python manage.py runserver 0.0.0.0:8000
